# crontab
# crontab
0 0 */3 * * export $(cat /app/dotenv/.env.prod | xargs) && /usr/local/bin/python3 /app/main.py >> /var/log/cron.log 2>&1
