# crontab
# crontab
0 21 * * * export $(cat /app/dotenv/.env.test | xargs) && /usr/local/bin/python3 /app/main.py >> /var/log/cron.log 2>&1
